# systemd/hgnotifier.service
[Unit]
Description=HG Notifier

[Service]
ExecStart=/usr/bin/hgdaemon
