# systemd/snotify.service
[Unit]
Description=Snotify daemon

[Service]
ExecStart=/usr/bin/snotifyd

[Install]
WantedBy=multi-user.target
